# deploy/invoice-audit.service.example
[Unit]
Description=Invoice Audit Agent Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/InvoiceAuditAgent/backend
EnvironmentFile=/opt/InvoiceAuditAgent/.env
ExecStart=/opt/InvoiceAuditAgent/backend/.venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
